# qemu/fs/etc/profile.d/am5_env.sh
# Load accessmod ENV variables
# ----------------------------
# This file will be overwritten by provision.sh 
# (keep it here, empty,  as ref)
