# qemu/fs/etc/profile.d/am5_env.sh
#!/bin/bash 
#
echo "Updated by provision.sh"
#

